# crontab
40 14 * * * cd /app && . ./env.sh && ./reload.sh
