# crontab
40 12 */5 * * cd /app && . ./env.sh && ./refresh-trends.sh
